# autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=NetworkManager Applet
Exec=nm-applet
Hidden=true
